# supervisor.conf
[program:gp]
command=sh run.sh
directory=/wwwroot/scrapy-spiders ; 执行 command 之前，先切换到工作目录
user=root
autorestart=true
startretries=99999